# stage3/01-sys-tweaks/files/profile.d/frc_pkgconfig.sh
export PKG_CONFIG_PATH=/usr/local/frc/lib/pkgconfig
