# wlanpi1/02-setup-services/files/lib/systemd/system/iperf3.service
[Unit]
Description=iperf3 server
After=network-online.target

[Service]
Type=idle
ExecStart=/usr/bin/iperf3 -s -i1

[Install]
WantedBy=network-online.target
