# wlanpi/03-pipx/files/lib/systemd/system/profiler.service
[Unit]
Description=profiler service
After=multi-user.target

[Service]
Type=idle
ExecStart=/opt/wlanpi/pipx/bin/profiler
