# wlanpi/02-setup-services/files/lib/systemd/system/iperf2.service
[Unit]
Description=iperf2 server
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/iperf -s -i1
