# wlanpi1/05-setup-services/files/lib/systemd/system/iperf2-udp.service
[Unit]
Description=iperf2 UDP server
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/iperf -s -i1 -u -p 5002
