# wlanpi1/05-setup-services/files/lib/systemd/system/pcie.service
[Unit]
Description=Enable PCIe

[Service]
Type=oneshot
ExecStart=dtparam pcie=on
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
