# wlanpi1/01-config-files/files/etc/systemd/system/kismet.service.d/override.conf
[Service]
User=wlanpi
Group=kismet
ExecStopPost=+/opt/wlanpi-common/wlanpi-stop-monitor-mode.sh
