# wlanpi1/10-first-boot/files/lib/systemd/system/wlanpi-first-boot.service
[Unit]
Description=Script executed only once at the first boot to finish image configurations

[Service]
ExecStart=/usr/bin/first-boot.sh

[Install]
WantedBy=multi-user.target
